# task_5/01-my-first-udev.rules


ACTION=="add", RUN+="/lib/udev/atclib.sh"

